# packer/shared/files/set-hostname.service
[Unit]
Description=Set Hostname and DNS

[Service]
ExecStart=/bin/bash /usr/bin/set_hostname.sh

[Install]
WantedBy=multi-user.target
